# roles/graphite/templates/supervisor-graphite.conf.j2
[program:gunicorn-graphite]
command=/opt/graphite/bin/gunicorn_django -u {{ nginx_user }} -g {{ nginx_user }} -b 127.0.0.1:8080 --log-file=/opt/graphite/storage/log/webapp/gunicorn.log /opt/graphite/webapp/graphite/settings.py
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=QUIT
user={{ nginx_user }}
 
[program:carbon-cache]
command=/opt/graphite/bin/python /opt/graphite/bin/carbon-cache.py --debug start
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=QUIT