# roles/statsd/templates/supervisor-statsd.conf.j2
[program:statsd]
command={{ node_path.stdout }} /opt/statsd/stats.js /opt/statsd/localConfig.js
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=QUIT
